# etc/cron.daily/firewall-tester
#!/bin/bash

/usr/sbin/firewall-tester -b
